# dependencies/crontab
# run hello.sh at midnight on the first of the month
0 0 1 * * /dependencies/generate_certs.sh
